# Setup/minecraft.service
[Service]
Type=simple
ExecStart=python3 ServerWrapper.py
WorkingDirectory=/home/minecraft/BYU_Servers/ServerControlScripts
User=minecraft
Restart=always
TimeoutStopSec=30
KillMode = mixed
[Install]
WantedBy=multi-user.target
