# Setup/minecraft.service
[Service]
Type=simple
ExecStart=python3 ServerWrapper.py
WorkingDirectory=/home/byumc/BYU_Servers/ServerControlScripts
User=byumc
Restart=always
TimeoutStopSec=120
KillMode = mixed
[Install]
WantedBy=multi-user.target
